# assets/supervisord.dockerd.conf
[program:dockerd]
command=dockerd
autostart=true
autorestart=true
stdout_logfile=/var/log/dockerd.log
redirect_stderr=true